# deploy/netmon-engine.service
[Unit]
Description=lnmp Monitoring Engine
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=netmon
Group=netmon
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW
WorkingDirectory=/opt/netmon/noop
EnvironmentFile=/etc/netmon/netmon.env
Environment=PYTHONPATH=/opt/netmon/noop/backend:/opt/netmon/noop
ExecStart=/opt/netmon/venv/bin/python -m monitoring.engine
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netmon-engine

[Install]
WantedBy=multi-user.target
